# systemd/rubbish-empty.timer
[Unit]
Description=Run `rubbish empty` daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
