# systemd/rubbish-empty.service
[Unit]
Description=Run `rubbish empty`

[Service]
Type=oneshot
ExecStart=rubbish empty 30 --no-confirm
