# systemd/trash-empty.service
[Unit]
Description=Run `trash empty`

[Service]
Type=oneshot
ExecStart=trash empty 30 --no-confirm
