# systemd/trash-empty.timer
[Unit]
Description=Run `trash empty` daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
